# DEBIAN/postinst
#!/bin/bash
systemctl daemon-reexec
systemctl daemon-reload
systemctl enable backuper.service
systemctl start backuper.service
